# rc.local
sudo python3 /home/ece2799team5/Documents/ECE\ 2799\ TEAM\ 5/finalrun.py